# roles/ha/templates/etc/cron.d/sitessync
* * * * * root /usr/sbin/csync2_sync -xv
